# resources/supervisord.conf
[supervisord]
nodaemon=true

[program:xltv]
command=/opt/xltv/server/bin/server.sh
